# deployment/crm-backend.service
[Unit]
Description=CRM Backend Application
After=network.target
Documentation=https://github.com/your-repo/crm-backend

[Service]
Type=simple
User=crm
Group=crm
WorkingDirectory=/opt/crm-backend
Environment="PATH=/opt/crm-backend/venv/bin"
Environment="PYTHONUNBUFFERED=1"

# Start command
ExecStart=/opt/crm-backend/venv/bin/python start_server.py

# Graceful shutdown with 30 second timeout
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Restart policy
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crm-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/crm-backend/data.db /opt/crm-backend/backups /opt/crm-backend/uploads

[Install]
WantedBy=multi-user.target
